# init/pirgb-web.service
[Unit]
Description=Web server (with ui) for control the pirgb-server
After=network.target

[Service]
ExecStart=%h/.local/bin/pirgb-web

[Install]
WantedBy=multi-user.target
